# core/busybox/bbklogd.service
[Unit]
Description=BusyBox klogd

[Service]
Restart=always
ExecStart=/bin/busybox klogd -n

[Install]
WantedBy=syslog.target
