# core/busybox/bbsyslogd.service
[Unit]
Description=BusyBox syslogd

[Service]
Restart=always
ExecStart=/bin/busybox syslogd -n -C1024

[Install]
WantedBy=syslog.target
