# core/systemd-units/colorum.service
[Unit]
Description=colorum run

[Service]
Type=oneshot
WorkingDirectory=/etc/colorum
ExecStart=/usr/bin/mk

