# core/systemd-units/crond.service
[Unit]
Description=BusyBox crond

[Service]
Restart=always
ExecStart=/usr/bin/crond -f -l 5
