# core/systemd-units/klogd.service
[Unit]
Description=BusyBox klogd

[Service]
Restart=always
CapabilityBoundingSet=CAP_SYSLOG
ExecStart=/usr/bin/klogd -n
StandardOutput=null

