# core/systemd-units/klogd.service
[Unit]
Description=BusyBox klogd

[Service]
Restart=always
CapabilityBoundingSet=CAP_SYS_ADMIN
ExecStart=/sbin/klogd -n
StandardOutput=null

