# core/systemd-units/nscd.service
[Unit]
Description=nscd

[Service]
Type=forking
PIDFile=/run/nscd/nscd.pid
Restart=Always
ExecStart=/usr/bin/nscd
ExecStop=/usr/bin/nscd --shutdown

