# core/systemd-units/nscd.service
[Unit]
Description=nscd

[Service]
Type=forking
ExecStart=/usr/sbin/nscd

[Install]
WantedBy=multi-user.target
