# core/systemd-units/ntpd.service
[Unit]
Description=BusyBox ntpd

[Service]
Restart=always
RestartSec=60
EnvironmentFile=/etc/conf.d/net
CapabilityBoundingSet=CAP_SYS_TIME
ExecStart=-/usr/bin/ntpd -n -N $ntpdopts
StandardError=syslog

