# core/systemd-units/sshd.service
[Unit]
Description=SSH Secure Shell Service

[Service]
ExecStart=/usr/sbin/sshd -D -oKexAlgorithms=diffie-hellman-group-exchange-sha256

